# docker/legacy/xinetd.d/echo-dgram
service echo
{
    type        = INTERNAL
    id          = echo-dgram
    socket_type = dgram
    protocol    = udp
    wait        = yes
    user        = nobody
    disable     = no
    flags       = IPv4
}