# docker/legacy/xinetd.d/time-dgram
service time
{
    type        = INTERNAL
    id          = time-dgram
    socket_type = dgram
    protocol    = udp
    wait        = yes
    user        = nobody
    disable     = no
    flags       = IPv4
}